# jenkins-agent.service
[Unit]
Description=Jenkins Agent

[Service]
User=redgrist
WorkingDirectory=/home/jenkins
ExecStart=/bin/bash /usr/local/jenkins-service/start-agent.sh
Restart=always

[Install]
WantedBy=multi-user.target
